# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/sysconfig/httpd && exec /usr/sbin/httpd -DFOREGROUND"
